# deploy/selftend.service
[Unit]
Description=SelfTend Backend Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/selftend/backend
ExecStart=/opt/selftend/backend/selftend
Restart=always
RestartSec=5
Environment=GIN_MODE=release
EnvironmentFile=-/opt/selftend/backend/.env

[Install]
WantedBy=multi-user.target
